# distribution/src/deb/control-runtime/postinst
#!/bin/sh
update-rc.d -f openhab defaults
