# distribution/src/deb/control-runtime/prerm
#!/bin/sh
if [ "$1" = "remove" ]; then
    /etc/init.d/openhab stop
    update-rc.d -f openhab remove
fi

